# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/paho-mqtt-sn-gateway/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# Paho MQTT-SN Gateway finish script
#==============================================================================
bashio::log.info "[MQTT-SN] Gateway ended with exit code ${1} (signal ${2})..."
