# openthread_border_router/rootfs/etc/s6-overlay/s6-rc.d/paho-mqtt-sn-gateway/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Paho MQTT-SN Gateway start script
# ==============================================================================

set_conf_value () {
    local conf_file=$1
    local conf_key=$2
    local conf_val=$3
    sed -i "s/^#*${conf_key}=.*$/${conf_key}=${conf_val}/" $conf_file
}

. /etc/s6-overlay/scripts/otbr-agent-common

declare cnf_file
declare output_pipe
declare gateway_id
declare gateway_name
declare gateway_max_clients
declare broker_name
declare broker_port
declare gateway_ipv6_port
declare gateway_mcast_ipv6
declare gateway_mcast_hops
cnf_file="/etc/s6-overlay/s6-rc.d/paho-mqtt-sn-gateway/conf/gateway.conf"
output_pipe="/etc/s6-overlay/s6-rc.d/paho-mqtt-sn-gateway/scripts/output-pipe.sh"
gateway_id=$(bashio::config 'mqtt_sn_gateway_id')
gateway_name=$(bashio::config 'mqtt_sn_gateway_name')
gateway_max_clients=$(bashio::config 'mqtt_sn_gateway_max_clients')
broker_name=$(bashio::config 'mqtt_broker_name')
broker_port=$(bashio::config 'mqtt_broker_port')
broker_user=$(bashio::config 'mqtt_broker_user')
broker_pass=$(bashio::config 'mqtt_broker_pass')
gateway_ipv6_port=$(bashio::config 'mqtt_sn_ipv6_port')
gateway_mcast_ipv6=$(bashio::config 'mqtt_sn_mcast_ipv6')
gateway_mcast_hops=$(bashio::config 'mqtt_sn_mcast_hops')

if bashio::config.true 'mqtt_sn_mcast_ipv6_mesh_loc'; then
    mesh_local_prefix=$(ot-ctl prefix meshlocal | sed -n 1p | sed 's/Mesh Local Prefix: //' | awk -F '::' '{print $1}')
    gateway_mcast_ipv6="ff33:40:${mesh_local_prefix}::1"
    bashio::log.info "[MQTT-SN] Multicast IPv6 address set to mesh local: ${gateway_mcast_ipv6}"
fi

bashio::log.info "[MQTT-SN] Configuring gateway"
set_conf_value ${cnf_file} "GatewayID" ${gateway_id}
set_conf_value ${cnf_file} "GatewayName" ${gateway_name}
set_conf_value ${cnf_file} "MaxNumberOfClients" ${gateway_max_clients}
set_conf_value ${cnf_file} "BrokerName" ${broker_name}
set_conf_value ${cnf_file} "BrokerPortNo" ${broker_port}
set_conf_value ${cnf_file} "GatewayIPv6PortNo" ${gateway_ipv6_port}
set_conf_value ${cnf_file} "MulticastIPv6PortNo" ${gateway_ipv6_port}
set_conf_value ${cnf_file} "MulticastIPv6" ${gateway_mcast_ipv6}
set_conf_value ${cnf_file} "MulticastIPv6If" ${thread_if}
set_conf_value ${cnf_file} "MulticastHops" ${gateway_mcast_hops}

if [[ ! ${broker_user} == "" ]]; then
    set_conf_value ${cnf_file} "LoginID" ${broker_user}
fi

if [[ ! ${broker_pass} == "" ]]; then
    set_conf_value ${cnf_file} "Password" ${broker_pass}
fi

bashio::log.info "[MQTT-SN] Starting gateway"
exec stdbuf -oL /usr/sbin/MQTT-SNGateway -f ${cnf_file} \
    1> >(xargs -I {} ${output_pipe} "out" {}) \
    2> >(xargs -I {} ${output_pipe} "err" {})
